# com.user.blog-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.blog-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/chanwook/Documents/Claude/Projects/블로그/run_daily_auto.command</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/chanwook/Documents/Claude/Projects/블로그/logs/launchd_out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chanwook/Documents/Claude/Projects/블로그/logs/launchd_err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/chanwook/Documents/Claude/Projects/블로그</string>
</dict>
</plist>
